# systemd/daily.service
[Unit]
Description=Run Our Daily Script

[Service]
Type=oneshot
ExecStart=/home/trevorgrabham/bin/daily.sh
