# systemd/daily.timer
[Unit]
Description=Runs Our Daily Script

[Timer]
OnStartupSec=1min
AccuracySec=30s
Persistent=false

[Install]
WantedBy=default.target
